# TankPy/webcam.cron
*/1 * * * * sudo bash ~/fishtank/takepic.sh
